# deploy/systemd/agent-qq.service
[Unit]
Description=agent-qq QQ AI Agent
Documentation=file:__PROJECT_DIR__/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON_BIN__ __PROJECT_DIR__/bot.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
